# system_files/shared/etc/systemd/system/greetd.service.d/pty.conf
[Service]
# Ensure pty devices are accessible for terminal emulators like ghostty
# PrivateDevices must be disabled to access the real /dev/pts
PrivateDevices=no
